# FrameNapTime-User.timer
[Unit]
Description=FrameWakeup

[Timer]
Unit=FrameWakeup.service
OnCalendar=*-*-* 23:00
Persistent=false

[Install]
WantedBy=timers.target
  
